# deploy/langflow-chat.service
[Unit]
Description=AlfyAI Application
After=network.target

[Service]
Type=simple
User=alfyai
WorkingDirectory=/opt/alfyai
EnvironmentFile=/opt/alfyai/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
